# etc/systemd/system/update-rust-analyzer.timer
[Unit]
Description=Update rust-analyzer weekly

[Timer]
OnBootSec=15min
OnUnitActiveSec=1w

[Install]
WantedBy=timers.target
